# packages/debian/source/common/postrm
#!/bin/sh
set -e

systemctl --system daemon-reload >/dev/null || true

if [ "$1" = "remove" ]; then
	deb-systemd-helper mask serenedb.service >/dev/null || true
fi

if [ "$1" = "purge" ]; then
	deb-systemd-helper purge serenedb.service >/dev/null || true
	deb-systemd-helper unmask serenedb.service >/dev/null || true
	rm -rf /var/lib/serenedb
	rm -rf /var/log/serenedb
fi
