# packages/debian/source/common/postinst
#!/bin/sh
set -e

deb-systemd-helper unmask serenedb.service >/dev/null
deb-systemd-helper enable serenedb.service >/dev/null
systemctl --system daemon-reload >/dev/null
deb-systemd-invoke start serenedb.service >/dev/null || true
