# packages/debian/source/common/serenedb.service
[Unit]
Description=SereneDB database server
After=network.target

[Service]
Type=simple
User=serenedb
Group=serenedb

LimitNOFILE=131072
LimitNPROC=131072
TasksMax=131072

# Ensure correct directory ownership before start
# '+' prefix runs these as root (needed for chown/install -o)
ExecStartPre=+/usr/bin/env chown -R serenedb:serenedb /var/lib/serenedb /var/log/serenedb

ExecStart=/usr/bin/serened

TimeoutSec=3600
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
